# web-deploy.service
[Unit]
Description="Website Deployer Service"

[Service]
ExecStart=/usr/local/bin/web-deploy
